# hypr/hyprshell/systemd/hyprshell.service.d/override.conf
[Service]
# The packaged unit uses `exit 0` when Hyprland is absent, which allows
# ExecStart to run anyway. A false ExecCondition skips the service cleanly.
ExecStartPre=
ExecCondition=/bin/sh -c '[ "$XDG_CURRENT_DESKTOP" = "Hyprland" ]'
